# services/mpdscribble/run
#!/bin/sh
exec mpdscribble --no-daemon --conf /etc/mpdscribble.conf
